# apps/wicd/trunk/debian/postinst
#!/bin/sh
# postinst script for wicd
#
# see: dh_installdeb(1)

set -e

update-rc.d wicd defaults

exit 0


